# packaging/macos/com.suricatoos.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.suricatoos.agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/suricatoos-agent</string>
		<string>run</string>
		<string>--state</string>
		<string>/var/lib/suricatoos-agent</string>
		<string>--queue</string>
		<string>/var/lib/suricatoos-agent/queue</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/suricatoos-agent.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/suricatoos-agent.log</string>
</dict>
</plist>
